# postfix/supervisor/conf.d/dovecot.ini
[program:dovecot]
command=/usr/sbin/dovecot
process_name=%(program_name)s
numprocs=1
priority=998
startsecs=5
redirect_stderr=true
;autostart=false
