# postfix/supervisor/conf.d/postfix.ini
[program:postfix]
command=postfix start
process_name=%(program_name)s
numprocs=1
priority=999
startsecs=1
redirect_stderr=true