# postfix/supervisor/conf.d/rsyslog.ini
[program:rsyslog]
command=/usr/sbin/rsyslogd -n -c3
priority=9
autostart=false
