# supervisord.d/nginx.ini
[program:nginx]
command=/sbin/nginx -p /opt/lsm -c /opt/lsm/config/nginx/nginx.conf
process_name=%(program_name)s
numprocs=1
directory=/opt/lsm
umask=022
priority=999
autostart=true
startsecs=1
autorestart=true
startretries=3000
user=root
redirect_stderr=true
stdout_logfile=/opt/lsm/logs/%(program_name)s-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
